# watercooler.service
[Unit]
Description=Water Cooler Auto Speed Daemon
After=bluetooth.service
Wants=bluetooth.service

[Service]
Type=simple
ExecStart=__WATERCOOLER_DIR__/venv/bin/python3 __WATERCOOLER_DIR__/watercooler.py daemon
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Bluetooth needs dbus access
Environment=WATERCOOLER_CONF_DIR=__WATERCOOLER_DIR__
Environment=DBUS_SYSTEM_BUS_ADDRESS=unix:path=/run/dbus/system_bus_socket

[Install]
WantedBy=default.target
